# crontab
@reboot python3 /home/pi/check_pin19.py
